# cpupower.service
[Unit]
Description=CPU comfort-zone 
[Service]
Type=oneshot
ExecStart=/usr/bin/cpupower-set-range 400 2200 # Adjust to resemble Intel 2670QM base frequency
[Install]
WantedBy=multi-user.target
